# deploy/refactoring-factory.service
[Unit]
Description=Refactoring Factory
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/app/refactoring-factory
EnvironmentFile=/app/refactoring-factory/.env

ExecStart=/home/deploy/.local/bin/uv run uvicorn server:app \
    --host 127.0.0.1 \
    --port 8394 \
    --access-log

Restart=on-failure
RestartSec=5s
TimeoutStartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=refactoring-factory

[Install]
WantedBy=multi-user.target
